# deploy/price-scraper.service
[Unit]
Description=Market Fiyat Scraper
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/infdatabase
ExecStart=/root/.venv/bin/python -m pipeline.runner
EnvironmentFile=/root/infdatabase/.env
StandardOutput=journal
StandardError=journal
